# deploy/systemd/xdg-desktop-portal.service.d/desktop.conf
# Drop-in for the per-user xdg-desktop-portal service.
# The kiosk starts labwc bare (no XDG_CURRENT_DESKTOP), so the portal would not
# load its labwc config and would never expose the Camera interface Chromium needs.
# Setting it here makes the portal load /usr/share/xdg-desktop-portal/labwc-portals.conf
# and offer org.freedesktop.portal.Camera. See docs/6-camera.md step 2.
# Install to: ~/.config/systemd/user/xdg-desktop-portal.service.d/desktop.conf
[Service]
Environment=XDG_CURRENT_DESKTOP=labwc
